# recipes-app/voltcan/voltcan/voltcan
#!/bin/sh
### BEGIN INIT INFO
# Provides:          voltcan
# Required-Start:
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description:
# Description:
### END INIT INFO

PATH=/sbin:/bin:/usr/sbin:/usr/bin:/app/bin:/app
DESC=VoltCanService

set -e

source obc_common.sh

startService(){
    VoltCanService\
        -p 3016 \
        -c $OBC_OBCCONF_FILE \
        --log-file $OBC_LOG_PATH/canjm_volt \
        --update-rate 0 \
        --jm-mode voltage \
        --dump $OBC_LOG_PATH/voltcan_$(date +%Y%m%d_%H%M).dump \
        can0 >> $OBC_LOG_PATH/el_canjm_volt_$(date +%Y%m%d_%H%M).log 2>&1 &
        pid=$!
        echo "PID $pid"
}

case $1 in
  start)
        echo -n "Starting ${DESC}: "
        PID=$(pidof $DESC || true)
        if [[ -n "$PID" ]]; then
          echo "Process already running with PID $PID"
        else
          startService
        fi
  ;;
  stop)
        echo -n "Stopping ${DESC}: "
        PID=$(pidof $DESC || true)
        if [[ -n "$PID" ]]; then
          echo "PID $PID"
          kill $PID
        else
          echo "No running process found"
        fi
  ;;
  restart|force-reload)
        $0 stop
        $0 start
  ;;
  status)
        echo -n "Status ${DESC}: "
        PID=$(pidof $DESC || true)
        if [[ -n "$PID" ]]; then
          echo "RUNNING with PID ${PID}"
        else
          echo "NOT RUNNING"
        fi
   ;;
   *)
        echo "Usage: ${DESC} {start|stop|restart|force-reload|status}" >&2
        exit 1
        ;;
esac

exit 0

# vim:noet
